# overlays/flow-s-ebuilds/net-misc/FrankenBot/files/FrankenBot.initd
#!/sbin/openrc-run

JAVA_ARGS="\
-server"

depend() {
	need localmount net
}

name="FrankenBot"
description="An XMPP bot written in Scala using Smack"

command="/usr/bin/FrankenBot"
command_args=""

start_stop_daemon_args="-u ${FRANKENBOT_USER} \
			-N ${FRANKENBOT_NICE:-0} -w 2000"
pidfile="${FRANKENBOT_PIDFILE}"
retry="${FRANKENBOT_SIGTERM_TIMEOUT}"
