# overlays/khoverlay/sys-apps/system76-driver/files/system76-driver.openrc
#!/sbin/openrc-run

description="System76 airplane-mode hotkey and LED support"
command=/usr/bin/system76-daemon
command_background=true
pidfile=/run/system76-driver.pid
output_log=/var/log/system76-driver.out
error_log=/var/log/system76-driver.err

depend() {
    need dbus
}
